# share/amida-dashboard-frontend.service
[Unit]
Description=Amida Dashboard Frontend

[Service]
ExecStart=${BASEDIR}/frontend/run.sh
Restart=on-failure

[Install]
WantedBy=graphical.target
